# private_dot_config/private_systemd/user/update_devbox.service.tmpl
[Unit]
Description=Auto update devbox
After=local-fs.target
After=network-online.target
ConditionPathExists=/usr/bin/devbox

[Service]
Type=oneshot
ExecStart=/usr/bin/devbox version update